# build/ansible/roles/supervisord/files/ai-chat.ini
[program:ai-chat]
priority = 16
command = 
    /usr/bin/node
        /opt/ai-chat/mcp-client-chatbot/dist/src/server.js
user = pmm
directory = /opt/ai-chat/mcp-client-chatbot
autorestart = true
autostart = true
startretries = 10
startsecs = 1
stopsignal = TERM
stopwaitsecs = 30
stdout_logfile = /srv/logs/ai-chat.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 2
redirect_stderr = true
environment = 
    POSTGRES_URL="postgres://ai_chat_user:ai_chat_secure_password@localhost:5432/ai_chat",
    FILE_BASED_MCP_CONFIG="true",
    PORT="3001",
    NODE_ENV="production" 