# build/ansible/roles/supervisord/files/ai-chat.ini
[program:aichat-backend]
command=/usr/sbin/aichat-backend
args=-config /etc/aichat-backend/config.yaml
user=pmm
autostart=true
autorestart=true
startretries=1000
startsecs=1
redirect_stderr=true
stdout_logfile=/srv/logs/aichat-backend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
environment=
    AICHAT_PORT="3001",
    AICHAT_LLM_PROVIDER="openai",
    AICHAT_LLM_MODEL="gpt-4o-mini",
    AICHAT_MCP_SERVERS_FILE="/etc/aichat-backend/mcp-servers.json",
    AICHAT_LOG_LEVEL="info",
    AICHAT_DATABASE_URL="postgres://ai_chat_user:ai_chat_secure_password@127.0.0.1:5432/ai_chat?sslmode=disable",
    GIN_MODE="release" 