# build/ansible/roles/supervisord/files/grafana.ini
; It's a temporary file that will be replaced once pmm-managed 
; generates another `grafana.ini` after the first run

[program:grafana]
priority = 3
command =
    /usr/sbin/grafana server
        --homepath=/usr/share/grafana
        --config=/etc/grafana/grafana.ini
user = pmm
directory = /usr/share/grafana
autorestart = true
autostart = false
startretries = 10
startsecs = 1
stopsignal = TERM
stopwaitsecs = 300
stdout_logfile = /srv/logs/grafana.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 2
redirect_stderr = true
