# build/packages/config/pmm-agent.service
[Unit]
Description=pmm-agent
After=time-sync.target network.target

[Service]
User=pmm-agent
# We intentionally let systemd assign a default value to Group in order
# to prevent pmm-agent from being unable to access certain system resources.
# Group=root
Type=simple
ExecStart=/usr/sbin/pmm-agent --config-file=/usr/local/percona/pmm/config/pmm-agent.yaml
Restart=always
RestartSec=2s

[Install]
WantedBy=multi-user.target
