# dev/mysql-backup-compose/supervisord-mysql.ini
[program:mysqld-backup-dev]
priority = 3
command = /usr/local/bin/start-mysqld.sh
autorestart = true
autostart = true
startretries = 10
startsecs = 5
stopsignal = INT
stopwaitsecs = 60
environment = MYSQL_ROOT_PASSWORD="%(ENV_MYSQL_ROOT_PASSWORD)s"
stdout_logfile = /srv/logs/mysqld-backup-dev.log
stdout_logfile_maxbytes = 30MB
stdout_logfile_backups = 2
redirect_stderr = true

[program:mysql-backup-register]
priority = 20
command = /usr/local/bin/register-mysql.sh
autorestart = unexpected
autostart = true
startretries = 5
startsecs = 0
exitcodes = 0
environment = MYSQL_ROOT_PASSWORD="%(ENV_MYSQL_ROOT_PASSWORD)s"
stdout_logfile = /srv/logs/mysql-backup-register.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 2
redirect_stderr = true
