# .devcontainer/etc/supervisord.d/dbaas-controller.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:dbaas-controller]
priority = 6
command = /usr/sbin/dbaas-controller
user = pmm
autorestart = false
autostart = false
startretries = 10
startsecs = 1
stopsignal = TERM
stopwaitsecs = 300
stdout_logfile = /srv/logs/dbaas-controller.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 3
redirect_stderr = true
