# build/packer/ansible/roles/docker-setup/templates/pmm2.service
[Unit]
Description=PMM2 Docker container
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
ExecStartPre=-/usr/bin/docker exec pmm-server stop
ExecStartPre=-/usr/bin/docker rm pmm-server
ExecStart=/usr/bin/docker run --volume /srv/:/srv/ --rm --name pmm-server \
    --net host \
    {{ pmm_server_image_name }}

[Install]
WantedBy=default.target
