# build/ansible/roles/victorialogs/files/victorialogs.ini
[program:victoria-logs]
priority = 20
command = /usr/local/bin/victoria-logs 
    -storageDataPath=/srv/victorialogs/data
    -retentionPeriod=30d
    -httpListenAddr=0.0.0.0:9428
    -envflag.enable
    -envflag.prefix=VL_
user = pmm
autorestart = true
autostart = true
startretries = 1000
startsecs = 3
stopsignal = TERM
stopwaitsecs = 10
stdout_logfile = /srv/logs/victorialogs.log
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 3
redirect_stderr = true


[program:promtail]
priority = 21
command = /usr/bin/promtail -config.file /srv/victorialogs/promtail.yml
user = pmm
autorestart = true
autostart = true
startretries = 1000
startsecs = 3
stopsignal = TERM
stopwaitsecs = 10
stdout_logfile = /srv/logs/promtail.log
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 3
redirect_stderr = true
