# build/packer/ansible/roles/cloud-node/files/banner.service
[Unit]
Description=URL banner service
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/opt/show-url

[Install]
WantedBy=multi-user.target
