# build/packages/deb/prerm
#
# This script is executed in the pre-remove phase
#
#       $1=remove    : indicates a removal
#       $1=upgrade   : indicates an upgrade


set -e

case "$1" in

    remove)
        STOP_REQUIRED=true
        REMOVE_SERVICE=true
    ;;
    upgrade)
        if [ "$RESTART_ON_UPGRADE" = "true" ]; then
            STOP_REQUIRED=true
        fi
    ;;
    deconfigure|failed-upgrade)
    ;;

    *)
        echo "pre remove script called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# Stops the service
if [ "$STOP_REQUIRED" = "true" ]; then
    echo -n "Stopping pmm-agent service..."
    if command -v systemctl >/dev/null; then
        systemctl --no-reload stop pmm-agent.service || true
    fi
    echo " OK"
fi

if [ "$REMOVE_SERVICE" = "true" ]; then
    if command -v systemctl >/dev/null; then
        systemctl disable pmm-agent.service > /dev/null 2>&1 || true
    fi
fi

@scripts.footer@
