# testdata/supervisord.d/prometheus.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:prometheus]
priority = 7
command =
	/usr/sbin/prometheus
		--config.file=/etc/prometheus.yml
		--query.max-concurrency=30
		--storage.tsdb.path=/srv/prometheus/data
		--storage.tsdb.retention.time=30d
		--storage.tsdb.wal-compression
		--web.console.libraries=/usr/share/prometheus/console_libraries
		--web.console.templates=/usr/share/prometheus/consoles
		--web.enable-admin-api
		--web.enable-lifecycle
		--web.external-url=http://localhost:9090/prometheus/
		--web.listen-address=:9090
user = pmm
autorestart = true
autostart = true
startretries = 3
startsecs = 1
stopsignal = TERM
stopwaitsecs = 300
stdout_logfile = /srv/logs/prometheus.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 3
redirect_stderr = true
