# build/packer/ansible/roles/init-admin-password/files/init-admin-password.service
[Unit]
Description=init-admin-password
Wants=network-online.target
After=network-online.target
After=nss-user-lookup.target nss-lookup.target
After=time-sync.target
After=pmm-server.target

[Service]
Restart=on-failure
RestartSec=20
RemainAfterExit=true
Type=oneshot

ExecStart=/home/admin/init-admin-password

[Install]
WantedBy=default.target
