# managed/testdata/supervisord.d/grafana.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:grafana]
priority = 3
command =
    /usr/sbin/grafana-server
        --homepath=/usr/share/grafana
        --config=/etc/grafana/grafana.ini
        cfg:default.paths.data=/srv/grafana
        cfg:default.paths.plugins=/srv/grafana/plugins
        cfg:default.paths.logs=/srv/logs
        cfg:default.log.mode=console
        cfg:default.log.console.format=console
        cfg:default.server.root_url="https://%%(domain)s/graph"
environment =
    POSTGRES_ADDR="",
    POSTGRES_DBNAME="",
    POSTGRES_USERNAME="",
    POSTGRES_DBPASSWORD="",
    POSTGRES_SSL_MODE="disable",
    POSTGRES_SSL_KEY_PATH="",
    POSTGRES_SSL_CERT_PATH="",
    PERCONA_TEST_PMM_CLICKHOUSE_DATASOURCE_ADDR="127.0.0.1:8123",
user = grafana
directory = /usr/share/grafana
autorestart = true
autostart = true
startretries = 10
startsecs = 1
stopsignal = TERM
stopwaitsecs = 300
stdout_logfile = /srv/logs/grafana.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 2
redirect_stderr = true
