# managed/testdata/supervisord.d/grafana.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:grafana]
priority = 3
command =
    /usr/sbin/grafana-server
        --homepath=/usr/share/grafana
        --config=/etc/grafana/grafana.ini
        cfg:default.paths.data=/srv/grafana
        cfg:default.paths.plugins=/srv/grafana/plugins
        cfg:default.paths.logs=/srv/logs
        cfg:default.log.mode=console
        cfg:default.log.console.format=console
        cfg:default.server.root_url="https://%%(domain)s/graph"
environment =
    POSTGRES_ADDR="127.0.0.1:5432",
    POSTGRES_DBNAME="postgres",
    POSTGRES_USERNAME="db_username",
    POSTGRES_DBPASSWORD="db_password",
    PERCONA_TEST_POSTGRES_SSL_MODE="verify",
    PERCONA_TEST_POSTGRES_SSL_CA_PATH="path-to-CA-cert",
    PERCONA_TEST_POSTGRES_SSL_KEY_PATH="path-to-key",
    PERCONA_TEST_POSTGRES_SSL_CERT_PATH="path-to-cert",
    GF_DATABASE_HOST="127.0.0.1:5432",
    GF_DATABASE_USER="db_username",
    GF_DATABASE_PASSWORD="db_password",
    GF_DATABASE_SSL_MODE="verify",
    GF_DATABASE_CA_CERT_PATH="path-to-CA-cert",
    GF_DATABASE_CLIENT_KEY_PATH="path-to-key",
    GF_DATABASE_CLIENT_CERT_PATH="path-to-cert",
    PERCONA_TEST_PMM_CLICKHOUSE_DATASOURCE_ADDR="127.0.0.1:8123",
user = grafana
directory = /usr/share/grafana
autorestart = true
autostart = true
startretries = 10
startsecs = 1
stopsignal = TERM
stopwaitsecs = 300
stdout_logfile = /srv/logs/grafana.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 2
redirect_stderr = true
