# managed/testdata/supervisord.d/grafana.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:grafana]
priority = 3
command =
    /usr/sbin/grafana server
        --homepath=/usr/share/grafana
        --config=/etc/grafana/grafana.ini
        cfg:default.server.domain="192.168.0.42:8443"
environment =
    PMM_POSTGRES_ADDR="127.0.0.1:5432",
    PMM_POSTGRES_DBNAME="postgres",
    PMM_POSTGRES_USERNAME="db_username",
    PMM_POSTGRES_DBPASSWORD="db_password",
    PMM_POSTGRES_SSL_MODE="verify",
    PMM_POSTGRES_SSL_CA_PATH="path-to-CA-cert",
    PMM_POSTGRES_SSL_KEY_PATH="path-to-key",
    PMM_POSTGRES_SSL_CERT_PATH="path-to-cert",
    PMM_CLICKHOUSE_HOST="127.0.0.1",
    PMM_CLICKHOUSE_PORT="9000",
    PMM_CLICKHOUSE_USER="default",
    PMM_CLICKHOUSE_PASSWORD="clickhouse",
directory = /usr/share/grafana
autorestart = true
autostart = true
startretries = 10
startsecs = 1
stopsignal = TERM
stopwaitsecs = 300
stdout_logfile = /srv/logs/grafana.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 2
redirect_stderr = true
