# managed/testdata/supervisord.d/nomad-server.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:nomad-server]
priority = 5
command = /usr/local/percona/pmm/tools/nomad agent -config /srv/nomad/nomad-server-192.168.0.42:8443.hcl
autorestart = false
autostart = false
startretries = 10
startsecs = 1
stopsignal = INT
stopwaitsecs = 300
stdout_logfile = /srv/logs/nomad-server.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 3
redirect_stderr = true
