# managed/testdata/supervisord.d/otel-collector.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:otel-collector]
priority = 6
command = /usr/local/percona/pmm/tools/otelcol-contrib --config=/srv/otelcol/config.yaml
autorestart = true
autostart = false
startretries = 10
startsecs = 1
stopsignal = INT
stopwaitsecs = 60
stdout_logfile = /srv/logs/otel-collector.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 3
redirect_stderr = true
