# managed/testdata/supervisord.d/prometheus.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:prometheus]
command = /bin/echo Prometheus is substituted by VictoriaMetrics
user = pmm
autorestart = false
autostart = false
startretries = 10
startsecs = 1
stopsignal = TERM
stopwaitsecs = 300
stdout_logfile = /srv/logs/prometheus.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 3
redirect_stderr = true
