# managed/testdata/supervisord.d/qan-api2.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:qan-api2]
priority = 13
command =
	/usr/sbin/percona-qan-api2
		--data-retention=30
environment =
	PMM_CLICKHOUSE_ADDR="127.0.0.1:9000",
	PMM_CLICKHOUSE_DATABASE="pmm",


autorestart = true
autostart = true
startretries = 1000
startsecs = 1
stopsignal = TERM
stopwaitsecs = 10
stdout_logfile = /srv/logs/qan-api2.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 3
redirect_stderr = true
