# managed/testdata/supervisord.d/victoriametrics_envvars.ini
; Managed by pmm-managed. DO NOT EDIT.

[program:victoriametrics]
priority = 7
command =
	/usr/sbin/victoriametrics
		--promscrape.config=/etc/victoriametrics-promscrape.yml
		--retentionPeriod=30d
		--storageDataPath=/srv/victoriametrics/data
		--httpListenAddr=127.0.0.1:9090
		--search.disableCache=true
		--search.maxQueryLen=2MB
		--search.latencyOffset=10s
		--search.maxUniqueTimeseries=500000000
		--search.maxSamplesPerQuery=1600000000
		--search.maxQueueDuration=100s
		--search.logSlowQueryDuration=300s
		--search.maxQueryDuration=9s
		--promscrape.streamParse=false
		--prometheusDataPath=/srv/prometheus/data
		--http.pathPrefix=/prometheus
		--envflag.enable
		--envflag.prefix=VM_
user = pmm
autorestart = true
autostart = true
startretries = 10
startsecs = 1
stopsignal = INT
stopwaitsecs = 300
stdout_logfile = /srv/logs/victoriametrics.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 3
redirect_stderr = true
