# build/ansible/roles/ai-chat/files/ai-chat.ini
[program:aichat-backend]
command=/usr/bin/aichat-backend
args=--env-only
directory=/etc/aichat-backend
user=pmm
autostart=true
autorestart=true
startretries=1000
startsecs=1
redirect_stderr=true
stdout_logfile=/srv/logs/aichat-backend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=2
environment=
    AICHAT_PORT="3001",
    AICHAT_LLM_PROVIDER="openai",
    AICHAT_LLM_MODEL="gpt-4o-mini",
    AICHAT_MCP_SERVERS_FILE="/etc/aichat-backend/mcp-servers.json",
    AICHAT_LOG_LEVEL="info",
    GIN_MODE="release" 